# deploy/provisioning/templates/the-tale-supervisor.conf

{% macro program(name, priority) %}
[program:{{name}}]
command=/home/the-tale/current/venv/bin/django-admin dext_amqp_worker -w {{name}} --settings the_tale.settings
priority={{priority}}
autostart=false
autorestart=false
startsecs=1
startretries=0
stopsignal=TERM
# supervisor should never kill game workers
stopwaitsecs=1000000000
user=the-tale
group=the-tale
redirect_stderr=true
stdout_logfile=/var/log/the-tale/{{name}}.log
environment=HOME="/home/the-tale",PATH="/home/the-tale/current/venv/bin/"
directory=/home/the-tale/current/
{% endmacro %}

# start first, stop last
{{ program(name='bank_processor', priority=1)}}
{{ program(name='xsolla_banker', priority=2)}}

{{ program(name='refrigerator', priority=3)}}

{{ program(name='message_sender', priority=4)}}

{{ program(name='accounts_manager', priority=5)}}
{{ program(name='items_manager', priority=5)}}
{{ program(name='portal_long_commands', priority=5)}}
{{ program(name='linguistics_manager', priority=5)}}
{{ program(name='market_manager', priority=5)}}

# after collections_items_manager
{{ program(name='achievements_manager', priority=6)}}

{{ program(name='turns_loop', priority=8)}}
{{ program(name='logic_1', priority=8)}}
{{ program(name='logic_2', priority=8)}}
{{ program(name='highlevel', priority=8)}}
{{ program(name='game_long_commands', priority=8)}}
{{ program(name='pvp_balancer', priority=8)}}
{{ program(name='quests_generator', priority=8)}}

# starts after all other game workers
{{ program(name='supervisor', priority=9)}}

{{ program(name='registration', priority=10)}}

[group:portal]
programs=bank_processor, xsolla_banker, refrigerator, message_sender, accounts_manager, items_manager, portal_long_commands, linguistics_manager, market_manager, achievements_manager, registration
priority=2

[group:game]
programs=supervisor, logic_1, logic_2, highlevel, game_long_commands, pvp_balancer, quests_generator, turns_loop
priority=1