# deploy/provisioning/templates/the-tale-supervisor.conf

{% macro the_tale_worker(name, priority) %}
[program:{{name}}]
command=/home/the_tale/current/venv/bin/django-admin amqp_queues_worker -w {{name}} --settings the_tale.settings
priority={{priority}}
autostart=false
autorestart=false
startsecs=1
startretries=0
stopsignal=TERM
# supervisor should never kill game workers
stopwaitsecs=1000000000
user=the_tale
group=the_tale
redirect_stderr=true
stdout_logfile=/var/log/the_tale/{{name}}.log
environment=HOME="/home/the_tale",PATH="/home/the_tale/current/venv/bin/"
directory=/home/the_tale/current/
{% endmacro %}

{% macro tt_worker(user, name, priority) %}
[program:{{user}}]
command=/home/{{user}}/current/venv/bin/tt_service -s {{name}} -c /etc/{{user}}/config.json
priority={{priority}}
autostart=false
autorestart=false
startsecs=1
startretries=0
stopsignal=TERM
# supervisor should never kill game workers
stopwaitsecs=1000000000
user={{user}}
group={{user}}
redirect_stderr=true
stdout_logfile=/var/log/{{user}}/log.log
environment=HOME="/home/{{user}}",PATH="/home/{{user}}/current/venv/bin/"
directory=/home/{{user}}/current/
{% endmacro %}


{{ tt_worker(user='tt_diary', name='tt_diary', priority=0) }}
{{ tt_worker(user='tt_personal_messages', name='tt_personal_messages', priority=0) }}
{{ tt_worker(user='tt_storage', name='tt_storage', priority=0) }}
{{ tt_worker(user='tt_market', name='tt_market', priority=0) }}
{{ tt_worker(user='tt_soft_currency', name='tt_bank', priority=0) }}
{{ tt_worker(user='tt_players_timers', name='tt_timers', priority=0) }}
{{ tt_worker(user='tt_personal_impacts', name='tt_impacts', priority=0) }}
{{ tt_worker(user='tt_crowd_impacts', name='tt_impacts', priority=0) }}
{{ tt_worker(user='tt_job_impacts', name='tt_impacts', priority=0) }}
{{ tt_worker(user='tt_fame_impacts', name='tt_impacts', priority=0) }}
{{ tt_worker(user='tt_game_chronicle', name='tt_events_log', priority=0) }}
{{ tt_worker(user='tt_clans_chronicle', name='tt_events_log', priority=0) }}
{{ tt_worker(user='tt_money_spendings', name='tt_impacts', priority=0) }}
{{ tt_worker(user='tt_players_properties', name='tt_properties', priority=0) }}
{{ tt_worker(user='tt_clans_properties', name='tt_properties', priority=0) }}
{{ tt_worker(user='tt_matchmaker', name='tt_matchmaker', priority=0) }}
{{ tt_worker(user='tt_clans_currencies', name='tt_bank', priority=0) }}
{{ tt_worker(user='tt_emissaries_impacts', name='tt_impacts', priority=0) }}
{{ tt_worker(user='tt_places_effects', name='tt_effects', priority=0) }}
{{ tt_worker(user='tt_emissaries_events_currencies', name='tt_bank', priority=0) }}
{{ tt_worker(user='tt_emissaries_events_uniquer', name='tt_uniquer', priority=0) }}
{{ tt_worker(user='tt_discord', name='tt_discord', priority=0) }}
{{ tt_worker(user='tt_data_protector', name='tt_data_protector', priority=0) }}

# start first, stop last
{{ the_tale_worker(name='bank_processor', priority=1)}}
{{ the_tale_worker(name='xsolla_banker', priority=2)}}

{{ the_tale_worker(name='refrigerator', priority=3)}}

{{ the_tale_worker(name='message_sender', priority=4)}}

{{ the_tale_worker(name='items_manager', priority=5)}}
{{ the_tale_worker(name='linguistics_manager', priority=5)}}

# after collections_items_manager
{{ the_tale_worker(name='achievements_manager', priority=6)}}

{{ the_tale_worker(name='turns_loop', priority=8)}}
{{ the_tale_worker(name='logic_1', priority=8)}}
{{ the_tale_worker(name='logic_2', priority=8)}}
{{ the_tale_worker(name='quests_generator', priority=8)}}

# starts after all other game workers
{{ the_tale_worker(name='supervisor', priority=9)}}


[program:site]
command=/home/the_tale/current/venv/bin/gunicorn the_tale.wsgi:application -b 127.0.0.1:8000 -w 4
priority=3
autostart=false
autorestart=false
startsecs=1
startretries=0
stopsignal=QUIT
user=www-data
group=www-data
redirect_stderr=true
stdout_logfile=/var/log/the_tale/site.log
environment=HOME="/home/the_tale",PATH="/home/the_tale/current/venv/bin/"
directory=/home/the_tale/current/


[group:portal]
programs=bank_processor, xsolla_banker, refrigerator, message_sender, items_manager, linguistics_manager, achievements_manager
priority=2

[group:game]
programs=supervisor, logic_1, logic_2, quests_generator, turns_loop
priority=1


[group:service]
programs=tt_diary, tt_personal_messages, tt_storage, tt_market, tt_soft_currency, tt_players_timers, tt_personal_impacts, tt_crowd_impacts, tt_job_impacts, tt_fame_impacts, tt_game_chronicle, tt_clans_chronicle, tt_money_spendings, tt_players_properties, tt_clans_properties, tt_matchmaker, tt_clans_currencies, tt_emissaries_impacts, tt_places_effects, tt_emissaries_events_currencies, tt_emissaries_events_uniquer, tt_discord, tt_data_protector
priority=0
